# resources/patrick-web.conf
setuid patrick
setgid nogroup
start on runlevel [2345]
stop on runlevel [016]
respawn
exec bundle exec thin start -p 3000 -e production
